# homsai/rootfs/etc/services.d/homsai/run
#!/usr/bin/with-contenv bashio

java -jar /home/homsai-engine-0.0.1.jar --home_assistant.api_url=http://supervisor/core  --home_assistant.token="$SUPERVISOR_TOKEN"